# templogger.cron
00 12 * * * /usr/bin/python3 ~/TempLoggerBBC/BBCTempLogger.py
